# querynest.service
[Unit]
Description=QueryNest MongoDB Query Service
After=network.target

[Service]
Type=simple
User=querynest
WorkingDirectory=E:\QueryNest
ExecStart=C:\Python311\python.exe start.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=E:\QueryNest

[Install]
WantedBy=multi-user.target
